# debian/postrm
#!/bin/sh
update-rc.d -f nntpswitch remove >/dev/null 2>&1
